# generic/etc/rc.local
# Print the IP address
_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi

/usr/bin/tmux new-session -d '/root/bin/pingmonitor.pl' &

exit 0
